# linux/bash/git/mirror/systemd/git-mirror.timer
[Unit]
Description=Run Git Mirror Backup Daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
